# .config/systemd/user/diary-sync.timer
[Unit]
Description=Sync diary repository daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
