# .config/systemd/user/obsidian-sync.timer
[Unit]
Description=Sync obsidian repository hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
