# Library/LaunchAgents/com.user.diary-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.diary-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>cd /Users/neynt/diary && /usr/bin/git add -A && /usr/bin/git commit -m "Auto-commit $(date)" || true && /usr/bin/git pull --rebase && /usr/bin/git push</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/tmp/diary-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/diary-sync.err</string>
</dict>
</plist>
